# ansible/roles/zfs-backup-server/templates/pull-zfs-backups-after-wake.service
[Unit]
Description=Run zfs-backup-pull after wake from suspend or hibernate
After=suspend.target 

[Service]
User=root
Environment=DISPLAY=:0
ExecStart={{ zfsbackup_script_path }}

[Install]
WantedBy=suspend.target 

