# ansible/roles/zfs-backup-server/templates/pull-zfs-backups-after-wake.service
[Unit]
Description="Run zfs-backup-pull after wake from suspend or hibernate"
After=sleep.target hibernate.target suspend-then-hibernate.target hybrid-sleep.target 

[Service]
User={{ ansible_user }}
Group={{ ansible_user }}
Type=oneshot
ExecStart={{ zfsbackup_script_path }}
TimeoutSec=0

[Install]
WantedBy=sleep.target hibernate.target suspend-then-hibernate.target hybrid-sleep.target 

