# roles/backups-borg-client/templates/borg-backup.service.j2
[Unit]
Description=Borg offsite backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/borg-backup
StandardOutput=journal
StandardError=journal
