# ansible/roles/ansible-pull/templates/ansible-pull-after-wake.service
[Unit]
Description="Run ansible-pull after wake from suspend or hibernate"
After=sleep.target hybrid-sleep.target 

[Service]
Type=simple
RemainAfterExit=yes
ExecStart="{{ ansiblepull_workdir }}/{{ ansiblepull_script_name }}"
TimeoutSec=0

[Install]
WantedBy=sleep.target hybrid-sleep.target 

