# ansible/roles/wakeonlan/templates/wol.service.j2
[Unit]
Description=Enable Wake On Lan

[Service]
Type=oneshot
ExecStart = /sbin/ethtool -s {{wakeonlan_interface}} wol g

[Install]
WantedBy=basic.target
