# ansible/roles/composition-emailbackup/templates/emailbackup.timer
[Unit]
Description=Email Backup Sync Timer
Requires=emailbackup.service

[Timer]
# Run based on configuration (daily by default)
{% if emailbackup_sync_schedule == "hourly" %}
OnCalendar=hourly
{% elif emailbackup_sync_schedule == "weekly" %}
OnCalendar=weekly
{% else %}
OnCalendar=daily
{% endif %}
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
