# roles/system-zfs-policy/templates/zfs-snapshot-daily.timer.j2
[Unit]
Description=ZFS Daily Snapshot Timer

[Timer]
OnCalendar={{ system_zfs_policy_timer_daily }}
Persistent=true

[Install]
WantedBy=timers.target
