# roles/system-zfs-policy/templates/zfs-prune.timer.j2
[Unit]
Description=ZFS Snapshot Prune Timer

[Timer]
OnCalendar={{ system_zfs_policy_timer_prune }}
Persistent=true

[Install]
WantedBy=timers.target
