# roles/system-zfs-policy/templates/zfs-snapshot-monthly.service.j2
[Unit]
Description=ZFS Monthly Snapshot
After=zfs.target

[Service]
Type=oneshot
ExecStart={{ system_zfs_policy_script_path }}/zfs-snapshot --type monthly
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
