# roles/system-zfs-policy/templates/zfs-snapshot-monthly.timer.j2
[Unit]
Description=ZFS Monthly Snapshot Timer

[Timer]
OnCalendar={{ system_zfs_policy_timer_monthly }}
Persistent=true

[Install]
WantedBy=timers.target
