# ansible/roles/system-zfs-policy/templates/zfs-snapshot-hourly.timer.j2
[Unit]
Description=ZFS Hourly Snapshot Timer

[Timer]
OnCalendar={{ system_zfs_policy_timer_hourly }}
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
