# ansible/roles/system-zfs-policy/templates/zfs-snapshot-hourly.service.j2
[Unit]
Description=ZFS Hourly Snapshot
After=zfs.target

[Service]
Type=oneshot
ExecStart={{ system_zfs_policy_script_path }}/zfs-snapshot --type hourly
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
